# ace/launchd/ai.superace.cycle.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>ai.superace.cycle</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/andrewpilson/.openclaw/workspace/ace/launchd/run-ace-cycle.sh</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/andrewpilson/.openclaw/workspace</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>ACE_OPERATOR_CHANNEL</key>
    <string>telegram</string>
    <key>ACE_OPERATOR_TARGET</key>
    <string>7529788084</string>
    <key>ACE_DB_PATH</key>
    <string>/Users/andrewpilson/.openclaw/workspace/ace/state/ace.db</string>
    <key>ACE_BRIEFING_PATH</key>
    <string>/Users/andrewpilson/.openclaw/workspace/ace/state/ace_briefing.md</string>
    <key>ACE_TELEGRAM_TRANSPORT</key>
    <string>none</string>
    <key>ACE_NOTIFICATION_CHANNEL</key>
    <string>jace</string>
    <key>ACE_NOTIFICATION_TARGET</key>
    <string>7529788084</string>
  </dict>
  <key>RunAtLoad</key>
  <true/>
  <key>StartInterval</key>
  <integer>1800</integer>
  <key>StandardOutPath</key>
  <string>/Users/andrewpilson/.openclaw/logs/ace-cycle.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/andrewpilson/.openclaw/logs/ace-cycle.err.log</string>
</dict>
</plist>
